# crontab
# change path if needed
*/30 * * * * /root/run.sh
